# recipes/signalk-server/signalk-server/signalk.socket
[Unit]
Description=SignalK Server Socket
PartOf=signalk.service

[Socket]
ListenStream=3000

[Install]
WantedBy=sockets.target
